# fpm-consul/scripts/after-install.sh
#!/bin/bash

set -e

mkdir -p /var/consul
systemctl enable consul
